# docker/backend/crontab
# is equal to 16-30 CET, which is time when European Central Bank is updating exchange rates
30 15 * * * python /app/scripts/download_currency_exchange_rates.py