# docker/backend/crontab
* * * * * python /app/scripts/download_currency_exchange_rates.py