# docker/backend/crontab
# rates updated after midnight GMT
30 00 * * * python /app/scripts/download_currency_exchange_rates.py