# docker/backend/crontab
# is equal to 16-10 CET, which is time when European Central Bank is updating rates
10 15 * * * python /app/scripts/download_currency_exchange_rates.py